# automation/com.ibank.daily-gitsync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- StartCalendarInterval의 Hour/Minute는 UTC가 아니라 이 Mac의 시스템 시간대(로컬 시각) 기준. 시간대를 아시아/서울(한국)으로 맞춰 두면 평일 13:00 = KST(한국시간) 오후 1시. -->
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.ibank.daily-gitsync</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/jhcho/Documents/AX사업본부/아이뱅크_인턴/automation/daily-git-sync.sh</string>
	</array>
	<key>StartCalendarInterval</key>
	<array>
		<dict>
			<key>Weekday</key>
			<integer>1</integer>
			<key>Hour</key>
			<integer>13</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>2</integer>
			<key>Hour</key>
			<integer>13</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>3</integer>
			<key>Hour</key>
			<integer>13</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>4</integer>
			<key>Hour</key>
			<integer>13</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
		<dict>
			<key>Weekday</key>
			<integer>5</integer>
			<key>Hour</key>
			<integer>13</integer>
			<key>Minute</key>
			<integer>0</integer>
		</dict>
	</array>
	<key>StandardOutPath</key>
	<string>/tmp/ibank-daily-gitsync.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/ibank-daily-gitsync.err</string>
</dict>
</plist>
